# templates/default/php-fpm.service.erb
# Configured by Chef.
# Changes will be overwritten.
[Unit]
Description=The PHP <%= @version %> FastCGI Process Manager
Documentation=man:php-fpm<%= @version %>(8)
After=network.target

[Service]
Type=notify
PIDFile=<%= @pidfile %>
ExecStart=<%= @bin %> --nodaemonize --fpm-config <%= @conf_file %>
<% if @unix_socket -%>
ExecStartPost=-/usr/bin/update-alternatives --quiet --install /run/php/php-fpm.sock php-fpm.sock <%= @unix_socket %> 73
ExecStopPost=-/usr/bin/update-alternatives --quiet --remove php-fpm.sock <%= @unix_socket %>
<% end -%>
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
